# conf/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket


[program:anwen]
numprocs = 2
numprocs_start = 1
directory=/var/www/anwen         ; directory to cwd to before exec (def no cwd)
command=python3 /var/www/anwen/hello.py --port=800%(process_num)d
process_name = %(program_name)s%(process_num)d
stdout_logfile = /home/anwen/anwen.log
redirect_stderr=false          ; 将stderr重定向到stdout
stderr_logfile = /home/anwen/anwen_err.log
autostart=true                ; supervisord守护程序启动时自动启动tornado
autorestart=true              ; supervisord守护程序重启时自动重启tornado


[program:our-anniversary]
numprocs = 1
numprocs_start = 1
directory=/var/www/our-anniversary         ; directory to cwd to before exec (def no cwd)
command=python3 /var/www/our-anniversary/hello.py --port=600%(process_num)d
process_name = %(program_name)s%(process_num)d
stdout_logfile = /home/anwen/our-anniversary.log
redirect_stderr=false          ; 将stderr重定向到stdout
stderr_logfile = /home/anwen/our-anniversary_err.log
autostart=true                ; supervisord守护程序启动时自动启动tornado
autorestart=true              ; supervisord守护程序重启时自动重启tornado
